# init/supervisor.ini
[program:http-proxy]
directory = /path/to/http-proxy
command = /path/to/http-proxy/bin/http-proxy
process_name = %(program_name)s
numprocs = 1
autostart = false
startsecs = 5
startretries = 3
stdout_logfile = /var/log/http-proxy.stdout.log
stdout_logfile_maxbytes = 10MB
stderr_logfile = /var/log/http-proxy.stderr.log
stderr_logfile_maxbytes = 10MB
